# chirpstack/chirpstack-concentratord-target-dragino-lps8n/files/spidev-bind
#!/bin/sh /etc/rc.common

START=10
STOP=10
USE_PROCD=1

start() {
	modprobe spidev 2>/dev/null || true
	for dev in spi1.0 spi0.0; do
		if [ -d "/sys/bus/spi/devices/$dev" ]; then
			echo spidev > "/sys/bus/spi/devices/$dev/driver_override" 2>/dev/null || true
			[ -d /sys/bus/spi/drivers/spidev ] && echo "$dev" > /sys/bus/spi/drivers/spidev/bind 2>/dev/null || true
		fi
	done
	[ ! -e /dev/spidev0.0 ] && [ -e /dev/spidev1.0 ] && ln -sf /dev/spidev1.0 /dev/spidev0.0
}


